# web/server/src/sh/bin/kaa_box_startup
#!/bin/bash
# startup script
# chkconfig: 345 89 1
### BEGIN INIT INFO
# provides: kaa_box_startup
# Required-Start:
# Required-Stop:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Description: Kaa box startup
### END INIT INFO

NAME="kaa_box_startup"
LOG="/var/log/kaa_box_startup.log"
SCRIPTS_PATH="/usr/lib/kaa-sandbox/bin"

PATH="/sbin:/usr/sbin:/bin:/usr/bin"
export PATH

IS_LSB=false

if [ -f /lib/lsb/init-functions ]; then
    . /lib/lsb/init-functions
    IS_LSB=true
else
    . /etc/rc.d/init.d/functions
fi

function log_success_console() {

	if [ "$IS_LSB" = true ]; then
        log_success_msg $1
    else
        success && echo $1
    fi
    
}

. /lib/lsb/init-functions

case "$1" in
start)

	start_message=$(printf "%-50s\n" "Starting $NAME...")
    log_success_console "$start_message"
    echo $start_message >> $LOG

    echo "" > /etc/resolv.conf

    printf "%-50s\n" "Updating IP..." >> $LOG
    bash $SCRIPTS_PATH/update_hosts.sh >> $LOG 2>&1

    sed -i 1i"nameserver 8.8.8.8" /etc/resolv.conf

    printf "%-50s\n" "Starting sandbox..." >> $LOG

    echo 0 > /proc/sys/kernel/hung_task_timeout_secs

    stty -F /dev/tty1 -echo >> $LOG 2>&1
    # drop caches (minimize swapping)
    sync && sysctl -w vm.drop_caches=3 >> $LOG 2>&1
;;
stop)
	stop_message=$(printf "%-50s\n" "Stopping $NAME...")
    log_success_console "$stop_message"
    echo $stop_message >> $LOG
    date >> $LOG
;;

restart)
    $0 stop
    $0 start
;;

*)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
esac
